# pentaho/configs/pentaho.conf
[program:pentaho]
command=/opt/pentaho/biserver-ce/start-pentaho.sh
startsecs=10
numprocs=1
redirect_stderr=true
autorestart=unexpected
