# mongodb/configs/mongod.conf
[program:mongod]
command=mongod
redirect_stderr=true
autorestart=unexpected
