# tor/configs/tor.sv.conf
[program:tor]
command=service tor restart
redirect_stderr=true
autorestart=unexpected
