# p2pool/configs/litecoin.sv.conf
[program:litecoind]
command=bitcoind
redirect_stderr=true
autorestart=unexpected
