# p2pool/configs/litecoind.sv.conf
[program:litecoind]
command=litecoind
redirect_stderr=true
autorestart=unexpected
